# livecd/template/etc/systemd/system/installer.target
[Unit]
Description=Guided Installer
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes

[Install]
Alias=default.target
